# contrib/systemd-configs/transmission-daemon.service.d/local.conf
[Unit]
After=transmission-ns.service
Requires=transmission-ns.service

[Service]
NetworkNamespacePath=/run/netns/transmission
