# contrib/systemd-configs/openvpn-client@nl.service.d/local.conf
[Unit]
After=transmission-ns.service
Requires=transmission-ns.service

[Service]
CapabilityBoundingSet=cap_sys_admin
